# tailscale/etc/NetworkManager/dispatcher.d/50-tailscale-udp-gro
#!/bin/sh

STATE="${2:-}"

case "$STATE" in
  ""|up|dhcp4-change|dhcp6-change|connectivity-change|reapply)
    ;;
  *)
    exit 0
    ;;
esac

default_iface() {
  ip -o route get 8.8.8.8 2>/dev/null |
    awk '{ for (i = 1; i <= NF; i++) if ($i == "dev") { print $(i + 1); exit } }'
}

NETDEV="${1:-}"
DEFAULT_NETDEV="$(default_iface)"

if [ -n "$DEFAULT_NETDEV" ]; then
  NETDEV="$DEFAULT_NETDEV"
fi

[ -n "$NETDEV" ] || exit 0

case "$NETDEV" in
  lo|tailscale*|docker*|br-*|veth*|virbr*)
    exit 0
    ;;
esac

command -v ethtool >/dev/null 2>&1 || exit 0
ethtool -K "$NETDEV" rx-udp-gro-forwarding on rx-gro-list off >/dev/null 2>&1 || true
